# service_management/diamond.conf
[program:diamond]
command=/usr/local/bin/diamond_manager
directory=/home/oseary/
user=root
autostart=true
autorestart=true
