# service_management/oseary.conf
[program:oseary]
command=/usr/share/python/oseary/bin/oseary
user=oseary
autostart=true
autorestart=true
stderr_logfile=/var/log/oseary.err.log
stdout_logfile=/var/log/oseary.out.log
